# scripts/bluetooth-autoconnect.service
[Unit]
Description=Auto-connect to a specific Bluetooth speaker at boot
Requires=bluetooth.service
After=bluetooth.service
# 如果你使用的是 PulseAudio
After=pulseaudio.service
# 如果你使用的是 PipeWire
# After=pipewire.service pipewire-pulse.service
After=network-online.target

[Service]
Type=oneshot
User=liucy
Group=liucy
ExecStart=/usr/local/bin/bluetooth-autoconnect.sh
TimeoutStartSec=600
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
